# deploy/systemd/edge-live-engine.service
[Unit]
Description=Edge Live Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/edge
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/edge/.venv/bin/edge-live-engine --config /opt/edge/project/configs/golden_certification.yaml --snapshot_path /var/lib/edge/live_state.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
